# dist/iaroapp.service
[Unit]
Description=iaroapp
After=network.target

[Service]
Type=simple
User=http
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/usr/share/iaroapp
ExecStart=/usr/share/iaroapp/.venv/bin/python /usr/share/iaroapp/manage.py runserver 0.0.0.0:8000
RestartSec=30
Restart=on-failure

[Install]
WantedBy=multi-user.target
