# dist/iaroapp.service
[Unit]
Description=iaroapp
After=network.target

[Service]
Type=simple
User=http
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/usr/share/iaroapp
ExecStart=docker run -p 192.168.178.251:8000:8000 app-iaro
RestartSec=30
Restart=on-failure

[Install]
WantedBy=multi-user.target
